# web_file_exchanger_PKGBUILD/src/web_file_exchanger.service
[Unit]
Description=web_file_exchanger
After=network.target

[Service]
Restart=always
ExecStart=/usr/bin/web_file_exchanger_server
User=horst
EnvironmentFile=/etc/web_file_exchanger/server_conf_env.sh
Type=simple

[Install]
WantedBy=default.target
